# collections/ansible_collections/redhat/rhel_system_roles/tests/systemd/templates/foo.service.j2
{{ ansible_managed | comment }}
[Service]
ExecStart=/bin/sleep 3600
ExecReload=/bin/touch /tmp/foo-reloaded

[Install]
WantedBy=multi-user.target